# deploy/systemd/bn-api.service
[Unit]
Description=BrighterNepal Flask API (Gunicorn + gevent)
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
User=brighternepal
Group=brighternepal
WorkingDirectory=/opt/brighternepal/brighter-nepal-api
EnvironmentFile=/opt/brighternepal/.env
ExecStart=/opt/brighternepal/brighter-nepal-api/venv/bin/gunicorn -c gunicorn.conf.py run:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5
# Memory guard — kill if a worker leaks beyond 1.5 GB
MemoryMax=1536M
# OOM priority: sacrifice API before the database
OOMScoreAdjust=200

StandardOutput=journal
StandardError=journal
SyslogIdentifier=bn-api

[Install]
WantedBy=multi-user.target
